# com.antiransomware.tokendropper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antiransomware.tokendropper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Library/AntiRansomware/macos_token_dropper.py</string>
        <string>--key</string>
        <string>/Library/AntiRansomware/keys/ed25519_private.key</string>
        <string>--pid</string>
        <string>0</string>
        <string>--path</string>
        <string>/Users/Shared/Protected</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/antiransomware_tokendropper.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/antiransomware_tokendropper.err</string>
</dict>
</plist>
