# archive/misc/linux_broker.service
[Unit]
Description=Anti-Ransomware Netlink Broker (PQC-enforced)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/antiransomware
ExecStart=/usr/bin/python3 /opt/antiransomware/linux_broker.py --key /opt/antiransomware/keys/ed25519_private.key --ttl 300 --ops 4294967295 --quota 1000000000
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
